# mappings/Azure/FileIntegrityMonitoring.yaml
version: 1.0
ATT&CK version: 8.1
date: 1/21/2021
name: File Integrity Monitoring
author: 
contact: ctid@mitre-engenuity.org
organization: Center for Threat Informed Defense (CTID)
platform: Azure
tags: 
  - Azure Security Center
  - Azure Defender for Servers
  - File system
  - Registry
  - Windows
  - Linux
description: >
  File integrity monitoring (FIM), also known as change monitoring, examines 
  operating system files, Windows registries, application software, Linux 
  system files, and more, for changes that might indicate an attack.
  File Integrity Monitoring (FIM) informs you when changes occur to sensitive 
  areas in your resources, so you can investigate and address unauthorized 
  activity. 
techniques:
  - id: T1053
    name: Scheduled Task/Job
    scores:
      - category: Detect
        value: Partial
        sub-techniques:
          - id: T1053.001
            name: At (Linux)
          - id: T1053.002
            name: At (Windows)
          - id: T1053.003
            name: Cron
          - id: T1053.005
            name: Scheduled Tasks
          - id: T1053.006
            name: Systemd Timers
        comment: >
            These techniques result in Windows Registry or file system artifacts
            being created or modified which can be detected by this control.
            T1053.004 Launchd was excluded because this control does not support
            the MacOS operating system.
  - id: T1098
    name: Account Manipulation
    sub-techniques: 
      - id: T1098.004
        name: SSH Authorized Keys
        scores: 
          - function: Detect
            value: Partial
            comment: >
              The authorized_keys file can be monitored for modification to 
              detect an attempt by an adversary to maintain access to an account.
  - id: T1547
    name: Boot or Logon Autostart Execution
    scores:
      - category: Detect
        value: Partial
        sub-techniques:
          - id: T1547.001
            name: Registry Run Keys / Startup Folder
          - id: T1547.002
            name: Authentication Package
          - id: T1547.003
            name: Time Providers
          - id: T1547.004
            name: Winlogon Helper DLL
          - id: T1547.005
            name: Security Support Provider
          - id: T1547.006
            name: Kernel Modules and Extensions
          - id: T1547.008
            name: LSASS Driver
          - id: T1547.009
            name: Shortcut Modification
          - id: T1547.010
            name: Port Monitors
          - id: T1547.012
            name: Print Processors
        comment: >
            These techniques result in Windows Registry or file system artifacts
            being created or modified which can be detected by this control.
            T1053.007 Re-opened Applications and T1053.011 Plist Modification were excluded 
            because this control does not support the MacOS operating system.
  - id: T1037
    name: Boot or Logon Initialization Scripts
    scores:
      - category: Detect
        value: Partial
        sub-techniques:
          - id: T1037.001
            name: Logon Script (Windows)
          - id: T1037.003
            name: Network Logon Script
        comment: >
            These techniques result in Windows Registry or file system artifacts
            being created or modified which can be detected by this control.  Other sub-techniques
            were excluded because the were specific to the MacOS operating system.
    sub-techniques: 
      - id: T1037.001
        name: Logon Script (Windows)
        scores: 
          - function: Detect
            value: Partial
            comment: >
              This control can be used to monitor changes to the registry key 
              related to Logon scripts.  This monitoring may lead to false
              positives.
      - id: T1037.002
        name: Logon Script (Mac)
        scores: 
          - function: Detect
            value: Partial
            comment: >
              It appears that configuration for login hooks is stored in the 
              /var/root/Library/Preferences/com.apple.loginwindow preferences 
              file and therefore this control can be used to monitor this file 
              for changes.
      - id: T1037.004
        name: Rc.common
        scores: 
          - function: Detect
            value: Partial
            comment: >
              This control can be used to monitor changes to this file.
      - id: T1037.005
        name: Startup Items
        scores: 
          - function: Detect
            value: Significant
            comment: >
              This control can be used to monitor the /Library/StartupItems
              for changes. Changes to this folder should be infrequent especially
              since it is a deprecated mechanism.
comments: >
    Although this control is capable of reliably detecting changes to sensitive
    resources that can be indicative of an attack, it does not include any 
    intelligence that would minimize false positives and consequently most of 
    its detection score has been limited to Minimal or Partial.
references:
  - https://docs.microsoft.com/en-us/azure/security-center/security-center-file-integrity-monitoring