# mappings/Azure/FileIntegrityMonitoring.yaml
version: '1.0'
ATT&CK version: '8.1'
creation date: 1/21/2021
name: File Integrity Monitoring
contact: ctid@mitre-engenuity.org
organization: Center for Threat Informed Defense (CTID)
platform: Azure
tags:
  - Azure Security Center
  - Azure Defender for Servers
  - File system
  - Registry
  - Windows
  - Linux
description: >
  File integrity monitoring (FIM), also known as change monitoring, examines  operating system
  files, Windows registries, application software, Linux  system files, and more, for changes that
  might indicate an attack. File Integrity Monitoring (FIM) informs you when changes occur to
  sensitive  areas in your resources, so you can investigate and address unauthorized  activity. 
techniques:
  - id: T1053
    name: Scheduled Task/Job
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1053.001
            name: At (Linux)
          - id: T1053.002
            name: At (Windows)
          - id: T1053.003
            name: Cron
          - id: T1053.005
            name: Scheduled Tasks
          - id: T1053.006
            name: Systemd Timers
        scores:
          - category: Detect
            value: Significant
  - id: T1098
    name: Account Manipulation
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1098.004
            name: SSH Authorized Keys
        scores:
          - category: Detect
            value: Significant
  - id: T1547
    name: Boot or Logon Autostart Execution
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1547.001
            name: Registry Run Keys / Startup Folder
          - id: T1547.002
            name: Authentication Package
          - id: T1547.003
            name: Time Providers
          - id: T1547.004
            name: Winlogon Helper DLL
          - id: T1547.005
            name: Security Support Provider
          - id: T1547.006
            name: Kernel Modules and Extensions
          - id: T1547.008
            name: LSASS Driver
          - id: T1547.009
            name: Shortcut Modification
          - id: T1547.010
            name: Port Monitors
          - id: T1547.012
            name: Print Processors
        scores:
          - category: Detect
            value: Significant
  - id: T1037
    name: Boot or Logon Initialization Scripts
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1037.001
            name: Logon Script (Windows)
          - id: T1037.003
            name: Network Logon Script
        scores:
          - category: Detect
            value: Significant
  - id: T1543
    name: Create or Modify System Process
    technique-scores:
      - category: Detect
        value: Significant
    sub-techniques-scores:
      - sub-techniques:
          - id: T1543.002
            name: Systemd Service
          - id: T1037.003
            name: Windows Service
        scores:
          - category: Detect
            value: Significant
  - id: T1546
    name: Event Triggered Execution
    technique-scores:
      - category: Detect
        value: Partial
        comment: >
          The detection score for this technique was assessed as Partial because it doesn't
          detect some of the sub-techniques of this technique such as Windows Management
          Instrumentation Event Subscription (WMI) and Trap sub-techniques.  Additionally for 
          some sub-techniques, this control can be noisy.
    sub-techniques-scores:
      - sub-techniques:
          - id: T1546.001
            name: Change Default File Association
          - id: T1546.002
            name: Screensaver
          - id: T1546.004
            name: .bash_profile and .bashrc
          - id: T1546.007
            name: Netsh Helper DLL
          - id: T1546.008
            name: Accessibility Features
          - id: T1546.009
            name: AppCert DLLs
          - id: T1546.011
            name: Application Shimming
          - id: T1546.012
            name: Image File Execution Options Injection
          - id: T1546.013
            name: PowerShell Profile
        scores:
          - category: Detect
            value: Significant
      - sub-techniques:
          - id: T1546.010
            name: AppInit DLLs
          - id: T1546.015
            name: Component Object Model Hijacking
        scores:
          - category: Detect
            value: Minimal
            comment: >
              The detection score for this group of sub-techniques is assessed as Minimal
              due to the accuracy component of the score.  The registry keys which are
              modified as a result of these sub-techniques can change frequently or are 
              too numerous to monitor and therefore can result in significant amount of false positives.
  - id: T1574
    name: Hijack Execution Flow
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1574.006
            name: LD_PRELOAD
        scores:
          - category: Detect
            value: Partial
  - id: T1137
    name: Office Application Startup
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1137.002
            name: Office Test
        scores:
          - category: Detect
            value: Significant
  - id: T1548
    name: Abuse Elevation Control Mechanism
    technique-scores:
      - category: Detect
        value: Minimum
    sub-techniques-scores:
      - sub-techniques:
          - id: T1548.002
            name: Bypass User Account Control
        scores:
          - category: Detect
            value: Minimum
      - sub-techniques:
          - id: T1548.003
            name: Sudo and Sudo Caching
        scores:
          - category: Detect
            value: Significant
  - id: T1036
    name: Masquerading
    technique-scores:
      - category: Detect
        value: Minimal
    sub-techniques-scores:
      - sub-techniques:
          - id: T1036.005
            name: Match Legitimate Name or Location
        scores:
          - category: Detect
            value: Minimal
comments: >
    These techniques included in this mapping result in Windows Registry or file system artifacts
    being created or modified which can be detected by this control.  
    
    The detection score for most techniques included in this mapping was scored as Significant 
    and where there are exceptions, comments have been provided. This Significant score assessment 
    was due to the following factors:
    Coverage - (High) The control was able to detect most of the sub-techniques, references and
    procedure examples of the mapped techniques.
    Accuracy - (High) Although this control does not include built-in intelligence to minimize 
    the false positive rate, the specific artifacts generated by the techniques in this mapping 
    do not change frequently and therefore the potential for a high false-positive is reduced. 
    Temporal - (Medium) This control at worst scans for changes on an hourly basis.
references:
  - https://docs.microsoft.com/en-us/azure/security-center/security-center-file-integrity-monitoring